# docker/supervisor/conf.d/cron.conf
[program:cron]
command=/usr/sbin/crond -f
autostart=true
autorestart=true
user=root
stderr_logfile=/var/www/storage/logs/cron.err.log
stdout_logfile=/var/www/storage/logs/cron.out.log
